# software/cs/dotnet/SensorClock/sensorclock.service
[Unit]
Description=SensorClock

[Service]
Type=notify
User=pi
ExecStart=/home/pi/sensorclock/SensorClock

[Install]
WantedBy=multi-user.target
