# software/cs/dotnet/SensorClock/sensorclock.service
[Unit]
Description=SensorClock

[Service]
Type=notify
User=barts
ExecStart=/home/barts/sensorclock/SensorClock

[Install]
WantedBy=multi-user.target
